# modules/avalon/templates/avalon_delayed_job_init_script.erb
#!/bin/sh
# Startup script for Avalon delayed_job daemon on RedHat/CentOS
# chkconfig: 2345 95 55
# description: Avalon delayed_job processing daemon
# processname: delayed_job

RAILS_ROOT="/var/www/avalon/current"
RAILS_ENV="<%=rails_env%>"

echo "Avalon delayed job $1..."
su - avalon -c "cd $RAILS_ROOT && bundle exec rake RAILS_ENV=$RAILS_ENV delayed_job:$1 2>/dev/null"