# modules/avalon/templates/avalon_delayed_job_init_script.erb
<%#
Copyright 2011-2013, The Trustees of Indiana University and Northwestern
  University.  Licensed under the Apache License, Version 2.0 (the "License");
  you may not use this file except in compliance with the License.

You may obtain a copy of the License at

http://www.apache.org/licenses/LICENSE-2.0

Unless required by applicable law or agreed to in writing, software distributed 
  under the License is distributed on an "AS IS" BASIS, WITHOUT WARRANTIES OR
  CONDITIONS OF ANY KIND, either express or implied. See the License for the 
  specific language governing permissions and limitations under the License.
---  END LICENSE_HEADER BLOCK  ---
%>
#!/bin/sh
# Startup script for Avalon delayed_job daemon on RedHat/CentOS
# chkconfig: 2345 95 55
# description: Avalon delayed_job processing daemon
# processname: delayed_job

RAILS_ROOT="/var/www/avalon/current"
RAILS_ENV="<%=rails_env%>"

echo "Avalon delayed job $1..."
su - avalon -c "cd $RAILS_ROOT && bundle exec rake RAILS_ENV=$RAILS_ENV delayed_job:$1 2>/dev/null"
